# deploy/node/services/comfyui.service
[Unit]
Description=ComfyUI Image Generation Server
After=network.target

[Service]
Type=simple
WorkingDirectory=__HOME__/ComfyUI
ExecStart=__HOME__/comfyui-start.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
